# etc/systemd/moros.socket
[Unit]
Description=moros public
PartOf=moros.service

[Socket]
ListenStream=42069
ReusePort=true
NoDelay=true
